# target/linux/ramips/base-files/etc/uci-defaults/02_network
#!/bin/sh

. /lib/functions.sh
. /lib/ramips.sh
. /lib/functions/uci-defaults.sh
. /lib/functions/system.sh

ramips_setup_gnubee_interfaces()
{
	uci batch <<EOF
set network.lan='interface'
set network.lan.ifname='eth0.1 eth1'
set network.lan.force_link=1
set network.lan.type='bridge'
set network.lan.proto='static'
set network.lan.ipaddr='192.168.1.1'
set network.lan.netmask='255.255.255.0'
set network.lan.ip6assign='60'
EOF

	ucidef_set_interface_wan "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 6t"
	ucidef_add_switch_vlan "switch0" "2" "4 6t"
}

ramips_setup_interfaces()
{
	local board="$1"

	ucidef_set_interface_loopback

	case $board in
	gb-pc1)
		ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
		ucidef_add_switch "switch0" "1" "1"
		ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 6t"
		ucidef_add_switch_vlan "switch0" "2" "4 6t"
		;;

        gb-pc2)
                ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
                ucidef_add_switch "switch0" "1" "1"
                ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 6t"
                ucidef_add_switch_vlan "switch0" "2" "4 6t"
                ;;

	gnubee)
		ramips_setup_gnubee_interfaces
		;;
	*)
		ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"

		;;
	esac
}

ramips_setup_macs()
{
	local board="$1"
	local lan_mac=""
	local wan_mac=""

	case $board in
	br-6475nd)
		lan_mac=$(cat /sys/class/net/eth0/address)
		wan_mac=$(mtd_get_mac_binary devdata 7)
		;;

	w306r-v20)
		lan_mac=$(cat /sys/class/net/eth0/address)
		wan_mac=$(macaddr_add "$lan_mac" 5)
		;;

	tew-691gr)
		lan_mac=$(cat /sys/class/net/eth0/address)
		wan_mac=$(macaddr_add "$lan_mac" 3)
		;;

	tew-692gr)
		lan_mac=$(cat /sys/class/net/eth0/address)
		wan_mac=$(macaddr_add "$lan_mac" 4)
		;;

	x5 |\
	x8)
		lan_mac=$(cat /sys/class/net/eth0/address)
		lan_mac=$(macaddr_add "$lan_mac" -2)
		;;


	hlk-rm04 | \
	mpr-a2 | \
	y1s)
		lan_mac=$(cat /sys/class/net/eth0/address)
		lan_mac=$(macaddr_setbit_la "$lan_mac")
		wan_mac=$(macaddr_add "$lan_mac" 1)
		;;

	cy-swr1100 | \
	dir-645)
		lan_mac=$(mtd_get_mac_ascii nvram lanmac)
		wan_mac=$(mtd_get_mac_ascii nvram wanmac)
		;;

	wcr-150gn)
                wan_mac=$(mtd_get_mac_binary factory 40)
		;;

	rt-n56u)
		lan_mac=$(cat /sys/class/net/eth0/address)
		lan_mac=$(macaddr_setbit_la "$lan_mac")
		wan_mac=$(mtd_get_mac_binary factory 32772)
		;;

	carambola | \
	freestation5 | \
	wnce2001)
		wan_mac=$(mtd_get_mac_binary factory 46)
		;;
	ralink-soc)
		lan_mac=$(eth_mac r lan)
                wan_mac=$(eth_mac r wan)
                ;;

	*)
		lan_mac=$(cat /sys/class/net/eth0/address)
		wan_mac=$(macaddr_add "$lan_mac" 1)
		;;

	esac

	[ -n "$lan_mac" ] && ucidef_set_interface_macaddr lan $lan_mac
	[ -n "$wan_mac" ] && ucidef_set_interface_macaddr wan $wan_mac
}

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

board=$(ramips_board_name)

ramips_setup_interfaces $board
ramips_setup_macs $board

uci commit network

exit 0
